# .config/systemd/user/upgrade-all-distroboxes.timer
[Unit]
Description=Daily run of upgrade-all-distroboxes

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
